# build/server/rtl_433.service
[Unit]
Description=rtl_433 SDR Receiver Daemon
After=network-online.target

[Service]
ExecStart=/usr/local/bin/rtl_433 -C si -M time:tz:utc -F "mqtt://localhost:1883,user=rtl_433,pass=***PASSWORD***,retain=0,events=rtl_433[/model][/id]"
Restart=always

[Install]
WantedBy=multi-user.target
